# root/etc/services.d/cryptpad/run
#!/usr/bin/with-contenv bash

cd / || exit

exec \
	s6-setuidgid abc node server.js
